# roles/consul/templates/consul.service.j2
[Unit]
Description=consul
Requires=network.target
After=syslog.target network.target

[Service]
Type=forking
ExecStartPre=-/bin/rm -f /run/consul/consul.pidExecStart={{consul_home}}/start.sh  
#ExecStartPre={{consul_home}}/start.sh
User=consul
EnvironmentFile={{ consul_home }}/consul.env
PIDFile=/run/consul/consul.pid
TimeoutSec=300

[Install]
WantedBy=multi-user.target
